# notify/sysd/fail-dummy.service
[Unit]
Description=Intentional failure service for testing
OnFailure=notify-email@fail-dummy.service

[Service]
Type=oneshot
ExecStart=/bin/false

